# configs/kiosk.conf
start on desktop-start

stop on desktop-end

respawn

exec /opt/kiosk/start-app.sh
